# wgcal/Batch/wgcal.cron
#!/bin/sh
# ---------------------------------------------------------
# $Id: wgcal.cron,v 1.1 2005/07/05 03:04:54 marc Exp $
# Distribued under GPL 2 licence terms
# ---------------------------------------------------------
#* */10 * * *  root  /home/httpd/what/wsh.php --api=wgcal_alert  > /dev/null 2>&1

