# freedom/FREEDOM.cron.in
#!/bin/sh
# ---------------------------------------------------------
# $Id: FREEDOM.cron.in,v 1.2 2005/11/10 16:00:36 eric Exp $
# 
# Distribued under GPL 2 licence terms
# ---------------------------------------------------------

psql freedom anakeen -f @prefix@/API/freedom_realclean.sql > /dev/null
