# freedom/FREEDOM.cron
#!/bin/sh
# ---------------------------------------------------------
# $Id: FREEDOM.cron,v 1.4 2003/07/29 13:08:59 eric Exp $
# 
# Distribued under GPL 2 licence terms
# ---------------------------------------------------------
find /var/run/php/ -name 'sess_*' -exec /bin/rm {} \;
psql freedom anakeen -f /home/httpd/what/API/freedom_clean.sql > /dev/null
