# freedom/FREEDOM.cron
#!/bin/sh
# ---------------------------------------------------------
# $Id: FREEDOM.cron,v 1.2 2002/07/24 07:22:53 eric Exp $
# 
# Distribued under GPL 2 licence terms
# ---------------------------------------------------------
psql freedom anakeen -f /home/httpd/what/API/freedom_clean.sql > /dev/null
